# scripts/wifi-autostart.service
# 车(AP)Wi-Fi 自启动(标志位机制) — systemd 版(不依赖桌面登录)。
# 安装(按车板实际路径改 ExecStart / User):
#   sudo cp wifi-autostart.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now wifi-autostart.service
# 看日志: journalctl -u wifi-autostart -f   或  ~/wifi_logs/wifi_ap.log
# 注意: 开关局域网仍用 ./stop_wifi.sh true|false + sudo reboot(改标志位);
#       本 service 只负责开机把 autostart_wifi.sh 拉起来。
[Unit]
Description=kian_flycar Wi-Fi AP autostart (car)
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=orangepi
ExecStart=/home/orangepi/kian_flycar/scripts/autostart_wifi.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
